# etc/systemd/maxou-clock.service
[Service]
ExecStart=/home/pi/maxou-clock/bin/maxou-clock.py
User=pi

[Install]
WantedBy=multi-user.target
