# etc/systemd/maxou-clock.service
[Service]
ExecStart=/home/pi/maxclock.py
User=pi

[Install]
WantedBy=multi-user.target
